# bin/config/supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=/usr/sbin/nginx

[program:dart-server]
command=/usr/bin/dart /var/www/daviscru/bin/server.dart