# deployment/fleet/data.service
[Unit]
Description=DavisCru Data container
After=etcd.service
After=docker.service

[Service]
User=core
EnvironmentFile=/etc/environment
TimeoutStartSec=0
KillMode=none
ExecStart=/bin/bash -c '/usr/bin/docker start data || /usr/bin/docker run -v /data --name data busybox echo Data-only container'