# deployment/fleet/data-staging.service
[Unit]
Description=DavisCru Staging Data container
After=etcd.service
After=docker.service

[Service]
User=core
EnvironmentFile=/etc/environment
TimeoutStartSec=0
KillMode=none
ExecStart=/bin/bash -c '/usr/bin/docker start data-staging || /usr/bin/docker run -v /data --name data-staging busybox echo Data-only staging container'